# nanobot/rootfs/etc/s6-overlay/s6-rc.d/init-nanobot/run
#!/usr/bin/with-contenv bashio
exec /usr/bin/nanobot-init
